# whisper-server/com.zeroxone.whisper-tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Phase 6c hotfix: SSH reverse tunnel LaunchAgent.

  Owner installs this at ~/Library/LaunchAgents/com.zeroxone.whisper-tunnel.plist
  via setup-mac-sidecar.sh, which substitutes HOME at install time.
  Loaded with `launchctl bootstrap gui/$(id -u) <plist>`.

  Maintains a reverse-port-forward from the Mac mini's localhost:9000
  (whisper-server FastAPI bind) to VPS:9000. The VPS sshd `GatewayPorts
  yes` directive then makes that listener bind to all interfaces — most
  importantly the docker bridge `172.17.0.1:9000` — so the bot
  container reaches it via `host.docker.internal` (mapped to
  `host-gateway` in compose).

  autossh keeps the tunnel alive across NAT timeouts and Mac-side
  network blips. Combined with launchd KeepAlive (process crash
  respawn), the tunnel survives both transport hiccups and process
  death.

  VPS host is hardcoded to the only deployment target (193.233.87.118).
  Authentication is via the dedicated SSH key generated by setup-mac-
  sidecar.sh (`~/.ssh/whisper_tunnel`). The corresponding public key
  is added to VPS authorized_keys with `restrict,permitlisten="9000"`
  so the key is good for nothing else (port-forward only, no shell,
  no other forwards).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zeroxone.whisper-tunnel</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/autossh</string>
        <string>-M</string>
        <string>0</string>
        <string>-N</string>
        <string>-R</string>
        <string>9000:localhost:9000</string>
        <string>0xone@193.233.87.118</string>
        <string>-i</string>
        <string>__HOME__/.ssh/whisper_tunnel</string>
        <string>-o</string>
        <string>ServerAliveInterval=30</string>
        <string>-o</string>
        <string>ServerAliveCountMax=3</string>
        <string>-o</string>
        <string>ExitOnForwardFailure=yes</string>
        <string>-o</string>
        <string>StrictHostKeyChecking=accept-new</string>
        <string>-o</string>
        <string>UserKnownHostsFile=__HOME__/.ssh/known_hosts</string>
        <string>-o</string>
        <string>BatchMode=yes</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__HOME__/whisper-server</string>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- AUTOSSH_GATETIME=0 disables the "first connection must
             succeed within N sec" check, so launchd can take care of
             the initial backoff loop while the network comes up. -->
        <key>AUTOSSH_GATETIME</key>
        <string>0</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME__/whisper-server/logs/whisper-tunnel.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/whisper-server/logs/whisper-tunnel.err</string>

    <!-- Throttle restart attempts to avoid hammering VPS sshd if the
         key is misconfigured or the host is unreachable. 30s grace. -->
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
