# docs/ops/launchd.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Phase 5 template — NOT wired by default. Copy to
  ~/Library/LaunchAgents/com.agent2.0xone-assistant.plist (customise the
  paths!) and load with `launchctl bootstrap gui/$(id -u) ...`.

  Phase 8 will own the out-of-process discussion; today the daemon is
  in-process and lives inside one Python event loop.
-->
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agent2.0xone-assistant</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>--directory</string>
        <string>/Users/agent2/Documents/0xone-assistant</string>
        <string>python</string>
        <string>-m</string>
        <string>assistant</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/agent2/Documents/0xone-assistant</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ExitTimeOut</key>
    <integer>15</integer>

    <key>StandardOutPath</key>
    <string>/Users/agent2/Library/Logs/0xone-assistant/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/agent2/Library/Logs/0xone-assistant/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
